# recipes-rt-mag/files/etc/profile.d/GoToRT-MaG_Folder.sh
#!/bin/sh

cd ~/RT_MaG_ToolBox/RT_MaG_Program/

